# puppet/hpjdwpd/templates/jdwpd.init.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          jdwpd
# Required-Start:    $remote_fs $network
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: jdwpd
# Description:       This script starts and stops the jdwpd
### END INIT INFO


BASE=<%= @install_dir %>
NAME="jdwpd"
PYTHON=`which python`
PIDFILE="/var/run/jdwpd.pid"

. /lib/lsb/init-functions
killwaitpids() {
        PID="$@"
	if [ -z "$PID" ]; then return; fi
        kill -TERM $PID
        for i in `seq 1 30`; do
                ps -p $PID >/dev/null || break
                echo -n "."
                sleep 1
        done
}


cd $BASE || exit 1
case "$1" in

	start)
        	echo -n "Starting $NAME: "
		su -c "$PYTHON ${BASE}/jdwpd.py 1>>jdwpd-run.log 2>&1 & echo \$! > $PIDFILE" <%= @jdwpd_user %>
		iptables -t nat -I PREROUTING 1 -p tcp --dport 8000 -j REDIRECT --to-port <%= @jdwpd_port %>
	        echo "done"
	;;
	stop)
	        echo -n "Stopping $NAME: "
		if [ -f $PIDFILE ]; then
		        killwaitpids $(cat $PIDFILE)
		        rm $PIDFILE
		fi
	        killwaitpids $(pgrep -f ${BASE}/jdwpd.py)
		iptables -t nat -D PREROUTING -p tcp --dport 8000 -j REDIRECT --to-port <%= @jdwpd_port %>
	        echo "done"
	;;
	status)
		status_of_proc "${BASE}/jdwpd.py" "$NAME" && exit 0 || exit $?
	;;
	*)
		ps faxu | grep ${BASE}/jdwpd.py
		exit 1
	;;
esac

