# puppet/hpjdwpd/templates/warden_jdwpd.cron.erb
*/5 * * * * <%= @jdwpd_user %> cd <%= @install_dir %>; /usr/bin/python warden_sender_jdwpd.py 1>><%= @install_dir %>/warden_sender_jdwpd.log 2>&1
