# puppet/hptelnetd/templates/warden-telnetd.cron.erb
*/5 * * * * <%= @telnetd_user %> cd <%= @install_dir %>; /usr/bin/python warden3-telnetd-sender.py 1>><%= @install_dir %>/warden3-telnetd-sender.log 2>&1
