# puppet/hpglastopf/templates/warden-glastopf.cron.erb
*/5 * * * * <%= glastopf_user %> /usr/bin/python <%= install_dir %>/warden/glastopf-reporter.py 1>><%= install_dir %>/warden/glastopf-reporter.log 2>&1
