# puppet/hpglastopf/templates/warden-glastopf.cron.erb
*/5 * * * * <%= @glastopf_user %> /usr/bin/python <%= @install_dir %>/warden/warden3-glastopf-sender.py 1>><%= @install_dir %>/warden/glastopf-sender.log 2>&1
