# puppet/hpglastopf/templates/warden_glastopf.cron.erb
*/5 * * * * <%= @glastopf_user %> /usr/bin/python <%= @install_dir %>/warden/warden_sender_glastopf.py 1>><%= @install_dir %>/warden/warden_sender_glastopf.log 2>&1
