# puppet/hpdio/templates/dio.init.erb
#!/bin/bash
### BEGIN INIT INFO
# Provides:          dio
# Required-Start:    networking p0f
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Nahazuje dio
### END INIT INFO

NAME="dionaea"

BASE="<%= @install_dir %>"
OPTIONS=""
USER="<%= @dio_user %>"
GROUP="<%= @dio_user %>"
PIDFILE="$BASE/var/$NAME.pid"

. /lib/lsb/init-functions
killwaitpids() {
        PID="$@"
        if [ -z "$PID" ]; then return; fi
        kill -TERM $PID
        for i in `seq 1 30`; do
                ps -p $PID >/dev/null || break
                echo -n "."
                sleep 1
        done
}

ulimit -c unlimited
ulimit -n 50000

cd $BASE || exit 1
case "$1" in
	start)
        	echo -n "Starting $NAME: "
		$0 status 1>/dev/null 2>/dev/null
		if [ $? -eq 0 ]; then
			echo "already running"
			exit 0
		fi
		${BASE}/bin/dionaea -u $USER -g $GROUP -w ${BASE} -p $PIDFILE -D
		echo "done"
	;;
	stop)
        	echo -n "Stopping $NAME: "
                if [ -f $PIDFILE ]; then
                        killwaitpids $(cat $PIDFILE)
                        rm $PIDFILE
                fi
                killwaitpids $(pgrep -f ${BASE}/warden_ca_http.py)
                echo "done"
        ;;
	force-reload|restart)
        	$0 stop
	        sleep 10
	        $0 start
	;;
	log)
        	tail -f ${BASE}/var/log/dionaea*log
        ;;
	status)
		status_of_proc "${BASE}/bin/dionaea" "$NAME" && exit 0 || exit $?
	;;
	*)
        	echo "Usage: $0 {start|stop|restart|force-reload}"
	        ps faxu | grep "bin/$NAME"
        	exit 1
	;;
esac

exit 0
