# puppet/hpdio/templates/dio.init.erb
#!/bin/bash
### BEGIN INIT INFO
# Provides:          dio
# Required-Start:    networking p0f
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Nahazuje dio
### END INIT INFO

NAME="dionaea"

BASE="<%= install_dir %>"
OPTIONS=""
USER="<%= dio_user %>"
GROUP="<%= dio_user %>"
PIDFILE="$BASE/var/$NAME.pid"

cd $BASE
ulimit -c unlimited
ulimit -n 50000

case "$1" in
  start)
        echo -n "Starting $NAME: "
        ${BASE}/bin/dionaea -u $USER -g $GROUP -w ${BASE} -p $PIDFILE -D
        echo "done"
        ;;
  stop)
        echo -n "Stopping $NAME: "
        kill -TERM `cat $PIDFILE`
        echo "done"
        ;;
  force-reload|restart)
        $0 stop
        sleep 10
        $0 start
        ;;
  log)
        tail -f ${BASE}/var/log/dionaea*log
        ;;
  *)
        echo "Usage: $0 {start|stop|restart|force-reload}"
        ps faxu | grep "bin/$NAME"
        exit 1
esac

exit 0
