# puppet/hpdio/templates/warden_dio.cron.erb
*/5 * * * * <%= @dio_user %> /usr/bin/python <%= @install_dir %>/warden/warden_sender_dio.py 1>><%= @install_dir %>/warden/warden_sender_dio.log 2>&1
