# puppet/hpdio/templates/warden-dio.cron.erb
*/5 * * * * <%= dio_user %> /usr/bin/python <%= install_dir %>/warden/warden3-dio-sender.py 1>><%= install_dir %>/warden/warden-dio-sender.log 2>&1
