# puppet/hpcowrie/templates/warden-cowrie.cron.erb
*/5 * * * * <%= @cowrie_user %> /usr/bin/python <%= @install_dir %>/warden/warden3-cowrie-sender.py 1>><%= @install_dir %>/warden/warden-cowrie-sender.log 2>&1
