# puppet/hpcowrie/templates/warden_cowrie.cron.erb
*/5 * * * * <%= @cowrie_user %> /usr/bin/python <%= @install_dir %>/warden/warden_sender_cowrie.py 1>><%= @install_dir %>/warden/warden_sender_cowrie.log 2>&1
