# puppet/hpcowrie/templates/warden-cowrie.cron.erb
*/5 * * * * <%= cowrie_user %> /usr/bin/python <%= install_dir %>/warden/warden3-kippo-sender.py 1>><%= install_dir %>/warden/warden-kippo-sender.log 2>&1
*/5 * * * * <%= cowrie_user %> /usr/bin/python <%= install_dir %>/warden/cowrie-reporter.py 1>><%= install_dir %>/warden/cowrie-reporter.log 2>&1
