# puppet/hpelastichoney/templates/warden_elastichoney.cron.erb
*/5 * * * * <%= @elastichoney_user %> cd <%= @install_dir %>; /usr/bin/python warden_sender_elastichoney.py 1>><%= @install_dir %>/warden_sender_elastichoney.log 2>&1
