# puppet/hpelastichoney/templates/warden-elastichoney.cron.erb
*/5 * * * * <%= @elastichoney_user %> /usr/bin/python <%= @install_dir %>/elastichoney-reporter.py 1>><%= @install_dir %>/warden3-elastichoney-sender.log 2>&1
