# puppet/hpelastichoney/templates/warden-elastichoney.cron.erb
*/5 * * * * <%= elastichoney_user %> /usr/bin/python <%= install_dir %>/warden/elastichoney-reporter.py 1>><%= install_dir %>/warden/elastichoney-reporter.log 2>&1
