# puppet/hpconpot/templates/warden-conpot.cron.erb
*/5 * * * * <%= conpot_user %> /usr/bin/python <%= install_dir %>/warden/conpot-reporter.py 1>><%= install_dir %>/warden/conpot-reporter.log 2>&1
