# puppet/hpucho/templates/uchoudp.init.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          uchoudp
# Required-Start:    $remote_fs $network
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: uchoudp
# Description:       This script starts and stops the uchoudp
### END INIT INFO


BASE=<%= install_dir %>
NAME="uchoudp"
PYTHON=`which python`
PIDFILE="/var/run/uchoudp.pid"

. /puppet/metalib/lib.sh

cd $BASE || exit 1
ulimit -n 65535

killwaitpids() {
        PID="$@"
	if [ -z "$PID" ]; then return; fi
        kill -TERM $PID
        for i in `seq 1 30`; do
                ps -p $PID >/dev/null || break
                echo -n "."
                sleep 1
        done
}


case "$1" in

	start)
        	echo -n "Starting $NAME: "
		
		$PYTHON ${BASE}/uchoudp.py 1>uchoudp.log 2>&1 &
	        echo $! > $PIDFILE

	        echo "done"
	;;
	stop)
	        echo -n "Stopping $NAME: "
		if [ -f $PIDFILE ]; then
		        killwaitpids $(cat $PIDFILE)
		        rm $PIDFILE
		fi
	        killwaitpids $(pgrep -f ${BASE}/uchoudp.py)
	        echo "done"
	;;
	status)
		pgrep -f ${BASE}/uchoudp.py 1>/dev/null
                if [ $? -eq 0 ]; then
                        rreturn 0 "running"
                fi
		rreturn 1 "not running"
	;;
	*)
		ps faxu | grep ${BASE}/uchoudp.py
		exit 1
	;;
esac

