# puppet/hpucho/templates/warden-uchoweb.cron.erb
*/5 * * * * <%= @uchoweb_user %> /usr/bin/python <%= @install_dir %>/warden3-uchoweb-sender.py 1>><%= @install_dir %>/warden3-uchoweb-sender.log 2>&1
