# puppet/hpucho/templates/warden_uchoudp.cron.erb
*/5 * * * * <%= @uchoudp_user %> cd <%= @install_dir %>; /usr/bin/python warden_sender_uchoudp.py 1>><%= @install_dir %>/warden_sender_uchoudp.log 2>&1
