# puppet/hpucho/templates/warden-uchotcp.cron.erb
*/5 * * * * <%= @uchotcp_user %> cd <%= @install_dir %>; /usr/bin/python warden3-uchotcp-sender.py 1>><%= @install_dir %>/warden3-uchotcp-sender.log 2>&1
