# puppet/hpucho/templates/warden-uchoudp.cron.erb
*/5 * * * * <%= @uchoudp_user %> cd <%= @install_dir %>; /usr/bin/python warden3-uchoudp-sender.py 1>><%= @install_dir %>/warden3-uchoudp-sender.log 2>&1
